# infra/observability/configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=warn

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5

[program:loki]
command=/usr/local/bin/loki -config.file=/etc/loki/loki.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:tempo]
command=/usr/bin/tempo -config.file=/etc/tempo/tempo.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:prometheus]
command=/usr/local/bin/prometheus --config.file=/etc/prometheus/prometheus.yml --storage.tsdb.path=/data/prometheus --storage.tsdb.retention.time=30d --web.enable-lifecycle --web.enable-remote-write-receiver --web.listen-address=0.0.0.0:9090
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:otel-collector]
command=/usr/bin/otelcol-contrib --config=/etc/otelcol-contrib/config.yaml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

[program:grafana]
command=/usr/share/grafana/bin/grafana server --homepath=/usr/share/grafana --config=/etc/grafana/grafana.ini cfg:default.paths.data=/data/grafana cfg:default.paths.logs=/var/log/grafana cfg:default.paths.plugins=/data/grafana/plugins cfg:default.security.admin_user=%(ENV_GF_SECURITY_ADMIN_USER)s cfg:default.security.admin_password=%(ENV_GF_SECURITY_ADMIN_PASSWORD)s
autostart=true
autorestart=true
user=grafana
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40
environment=HOME="/usr/share/grafana"
